# .config/systemd/user/offlineimap@.timer
[Unit]
Description=Offlineimap timer unit

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=default.target
