# .config/systemd/user/offlineimap@.service
[Unit]
Description=offlineimap sync job

[Service]
Type=simple
ExecStart=/usr/bin/offlineimap -a %i

