# .config/systemd/user/xorg.target
[Unit]
Description=Xorg Server
Wants=xorg.service
After=xorg.service
After=display-setup.service
